# Debugging/libbpf-tools/resmon/resmon.service.in
[Unit]
Description=mlxsw resource monitor

[Service]
Type=notify
ExecStart=@BINDIR@/resmon start
ExecStop=@BINDIR@/resmon stop
