# Debugging/libbpf-tools/resmon/resmon-exporter.service.in
[Unit]
Description=mlxsw resource monitor exporter
Requisite=resmon.service
After=resmon.service

[Service]
ExecStart=@BINDIR@/resmon-exporter -l 0.0.0.0:9417
